# packages/obot/obot-100.tar.gz/obot-100/files/obotd.service
[Unit]
Description=OBOTD - 24/7 channel service
After=multi-user.target

[Service]
DynamicUser=True
StateDirectory=obotd
LogsDirectory=obotd
CacheDirectory=obotd
ExecStart=/usr/local/bin/obotd
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
